# contrib/systemd/yggdrasil.service
[Unit]
Description=yggdrasil chain control / reverse proxy
Documentation=https://github.com/bowerscd/yggdrasil
After=network-online.target
Wants=network-online.target

[Service]
# `notify-reload` (paired with sd_notify(READY=1) in yggdrasil) means
# `is-active` reports `active` only after the chain listener (if any),
# proxy supervisor, and control socket have all bound — not just after
# the process forked. The `-reload` half also wires `systemctl reload
# yggdrasil` (default action: SIGHUP) to an authoritative rule rescan:
# the daemon emits RELOADING=1 / MONOTONIC_USEC before the watcher
# drains the trigger, then READY=1 once the supervisor has reconciled
# the new set, so systemd knows the reload settled before returning to
# the operator.
Type=notify-reload
NotifyAccess=main
ExecStart=/usr/bin/yggdrasil run
ExecReload=/bin/kill -HUP $MAINPID
Restart=on-failure
RestartSec=2s

# Dedicated unprivileged user. Created by the package's sysusers.d
# entry (contrib/sysusers.d/yggdrasil.conf) on install.
User=yggdrasil
Group=yggdrasil

# yggdrasil only needs CAP_NET_BIND_SERVICE if any derived rule binds
# a port below 1024 (e.g. 0.0.0.0:443 for HTTPS rules). Keep both lines
# in lock-step — CapabilityBoundingSet is the upper bound, AmbientCapabilities
# is what's actually granted. If every rule binds >= 1024, override
# both to empty via `systemctl edit yggdrasil`.
AmbientCapabilities=CAP_NET_BIND_SERVICE
CapabilityBoundingSet=CAP_NET_BIND_SERVICE

# systemd manages the runtime directory: created on start with the right
# ownership (User=/Group=), wiped on stop. There is intentionally no
# package-managed state directory: gateways and mid-chain relays have no
# daemon state on disk. Terminals that store ACME material under
# /var/lib/yggdrasil should create that directory and add a
# ReadWritePaths= drop-in for it.
RuntimeDirectory=yggdrasil
RuntimeDirectoryMode=0750

# Hardening — see docs/security.md → Operational hardening.
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
PrivateTmp=true
PrivateDevices=true
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectControlGroups=true
RestrictSUIDSGID=true
RestrictRealtime=true
LockPersonality=true
MemoryDenyWriteExecute=true
SystemCallArchitectures=native
ReadOnlyPaths=/etc/yggdrasil

# Log handler — JSON to journald is the default.
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
